# wireguard-logd/etc/systemd/system/wglogd.service
[Unit]
Description=Logging daemon for WireGuard peer handshakes

[Service]
ExecStart=/usr/local/bin/wglogd

[Install]
WantedBy=multi-user.target
